# 1-server_unit/resource/es/etc/init/esgclog.conf
start on runlevel [345]
respawn
exec /personium/personium-esgclog/esGCLog.sh
