# 3-server_unit/resource/nfs/etc/cron.d/log-delete-cron
30 21 * * * root find -L /personium/personium-core/log -name '*.log*' -mtime +31 -exec /bin/rm -f {} \; > /dev/null 2>&1
30 21 * * * root find -L /personium/personium-engine/log/log -name '*.log*' -mtime +31 -exec /bin/rm -f {} \; > /dev/null 2>&1
30 21 * * * root find -L /personium/logback/log -name '*.log*' -mtime +31 -exec /bin/rm -f {} \; > /dev/null 2>&1
